# rootfs/standard/etc/update-motd.d/10-mynode-header
#!/bin/bash

THIS_SCRIPT="header"
MOTD_DISABLE=""

toilet -f standard -F metal myNode

printf '\nWelcome to myNode!\n'
uname -a
printf '\n'
