# rootfs/standard/etc/update-motd.d/10-mynode-header
#!/bin/bash

THIS_SCRIPT="header"
MOTD_DISABLE=""

toilet -f standard -F metal MyNode

printf '\nWelcome to MyNode!\n'
uname -a
printf '\n'
