# rootfs/free/etc/update-motd.d/10-mynode-header
#!/bin/bash

THIS_SCRIPT="header"
MOTD_DISABLE=""

toilet -f standard -F metal "myNode    Free"

printf '\nWelcome to myNode!\n'
uname -a
printf '\n'
