# rootfs/standard/etc/systemd/system/glances.service
# MyNode glances service (system monitoring)
# /etc/systemd/system/glances.service

[Unit]
Description=Glances
After=network.target

[Service]
ExecStart=/usr/local/bin/glances -w
User=root
Group=root
TimeoutSec=120
Restart=always
RestartSec=60
Nice=5

[Install]
WantedBy=multi-user.target